# user/redshift-gtk.service.d/after.conf
[Unit]
After=xresources.service geoclue-agent.service
