# user/redshift-gtk.service.d/xresources.conf
[Unit]
After=window-manager.target
After=xresources.service
