# user/app@.service
[Unit]
Description=Application %i
After=window-manager.target

[Service]
ExecStart=/usr/bin/env %i

[Install]
WantedBy=window-manager.target
