# user/app@.service
[Unit]
Description=Application %i
After=display-server.target
After=xresources.service

[Service]
ExecStart=/usr/bin/env %i

[Install]
WantedBy=display-server.target
