# user/awesome.service
[Unit]
Description=Awesome window manager
Documentation=man:awesome(1)
After=display-manager.target
Requires=display-manager.target

[Service]
ExecStart=/usr/bin/awesome
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
