# user/insync.service
[Unit]
Description=Insync
Documentation=man:insync(1)
PartOf=window-manager.target
After=window-manager.target
After=xresources.service

[Service]
Type=oneshot
ExecStart=/usr/bin/insync start
ExecStop=/usr/bin/insync quit
RemainAfterExit=yes

[Install]
WantedBy=window-manager.target
