# user/keymap.service
[Unit]
Description=User keymaps
Requires=window-manager.target

[Service]
ExecStart=/usr/bin/xmodmap %h/.Xmodmap
RemainAfterExit=yes

[Install]
WantedBy=window-manager.target
